# deep-agent/config/supervisor.conf
; LangGraph Deep Web Agent - Supervisor Configuration

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[include]
files = /etc/supervisor/conf.d/*.conf

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; Web Server (FastAPI)
[program:web]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4 --log-level info
directory=/app
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/web.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1",ENVIRONMENT="production"
stopsignal=INT
stopwaitsecs=10
stopasgroup=true
killasgroup=true

; Celery Worker
[program:worker]
command=celery -A app.celery worker --loglevel=info --concurrency=4
directory=/app
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1",ENVIRONMENT="production"
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

; Celery Beat Scheduler
[program:beat]
command=celery -A app.celery beat --loglevel=info --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/app
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/beat.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1",ENVIRONMENT="production"
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

; Nginx Web Server
[program:nginx]
command=nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

; Redis (if running locally instead of separate container)
[program:redis]
command=redis-server --requirepass ${REDIS_PASSWORD:-redis_password} --appendonly yes
user=appuser
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/redis.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

; Flower Monitoring (Celery)
[program:flower]
command=celery -A app.celery flower --port=5555
directory=/app
user=appuser
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/flower.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1",ENVIRONMENT="production"
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

; Event Handlers
[eventlistener:process_exit]
command=python /app/scripts/supervisor_event_handler.py
events=PROCESS_STATE_EXITED
buffer_size=100

[eventlistener:process_fatal]
command=python /app/scripts/supervisor_event_handler.py
events=PROCESS_STATE_FATAL
buffer_size=100

; Group definitions
[group:webapp]
programs=web,nginx

[group:workers]
programs=worker,beat

; Process environment variables
[program:web]
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1",ENVIRONMENT="production",DATABASE_URL="postgresql://deepagent_user:${POSTGRES_PASSWORD:-secure_password}@postgres:5432/deep_agent",REDIS_URL="redis://:${REDIS_PASSWORD:-redis_password}@redis:6379/0",SECRET_KEY="${SECRET_KEY:-your-secret-key-here}"

[program:worker]
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1",ENVIRONMENT="production",DATABASE_URL="postgresql://deepagent_user:${POSTGRES_PASSWORD:-secure_password}@postgres:5432/deep_agent",REDIS_URL="redis://:${REDIS_PASSWORD:-redis_password}@redis:6379/0",CELERY_BROKER_URL="pyamqp://deepagent_user:${RABBITMQ_PASSWORD:-rabbitmq_password}@rabbitmq:5672//",CELERY_RESULT_BACKEND="redis://:${REDIS_PASSWORD:-redis_password}@redis:6379/1",SECRET_KEY="${SECRET_KEY:-your-secret-key-here}"

[program:beat]
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1",ENVIRONMENT="production",DATABASE_URL="postgresql://deepagent_user:${POSTGRES_PASSWORD:-secure_password}@postgres:5432/deep_agent",REDIS_URL="redis://:${REDIS_PASSWORD:-redis_password}@redis:6379/0",CELERY_BROKER_URL="pyamqp://deepagent_user:${RABBITMQ_PASSWORD:-rabbitmq_password}@rabbitmq:5672//",CELERY_RESULT_BACKEND="redis://:${REDIS_PASSWORD:-redis_password}@redis:6379/1",SECRET_KEY="${SECRET_KEY:-your-secret-key-here}"

[program:flower]
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1",ENVIRONMENT="production",CELERY_BROKER_URL="pyamqp://deepagent_user:${RABBITMQ_PASSWORD:-rabbitmq_password}@rabbitmq:5672//",CELERY_RESULT_BACKEND="redis://:${REDIS_PASSWORD:-redis_password}@redis:6379/1"